# server/srcs/supervisord.conf
[supervisord]
nodaemon=true
[program:php7.3-fpm]
command=service php7.3-fpm start
[program:nginx]
command=service nginx start
[program:mysql]
command=service mysql start